# application/crontab
* * * * * /usr/local/bin/php -q ./CreateImages.php -e
